# docker/fulltclash.conf
[program:fulltclash]
command=python3 -u /app/main.py
directory=/app
redirect_stderr=true
stdout_logfile=/var/log/fulltclash.log