# docker/fulltclash.conf
[program:fulltclash]
command=python3 /app/main.py
directory=/app
stdout_logfile=/var/log/fulltclash.log